# python_server.service
[Unit]
Description=Python HTTP Server
After=network.target

[Service]
User=<your-user>
Group=<your-group>
WorkingDirectory=/path/to/your/script
ExecStart=/usr/bin/python3 /path/to/your/script/server.py
Restart=always

[Install]
WantedBy=multi-user.target
